# AdobeColdFusion/build/service/coldfusion9/finish
#!/bin/bash

set -e

/opt/coldfusion8/bin/coldfusion stop 2>&1

running=1

while [ $running = "1" ]
do
	status="$(ps auxw | grep coldfusion | grep autorestart | wc -l)"
	if [ "$status" != "1" ]
	then
		running=0
		break
	fi

	sleep 2
done

exit 1
